# ubl/infra/lab-256/api.service
[Unit]
Description=UBL API
[Service]
ExecStart=/usr/bin/ubl-api
