# ubl/infra/lab-256/tail.service
[Unit]
Description=UBL Tail SSE
[Service]
ExecStart=/usr/bin/ubl-tail
